# celery-worker.service
[Unit]
Description=Celery Worker for NxSocial
Documentation=https://docs.celeryq.dev
After=network.target redis.service postgresql.service

[Service]
Type=forking
User=www-data
Group=www-data

# ── IMPORTANT: Update these paths to match your server setup ──
WorkingDirectory=/var/www/nxsocial/Socialmedia
Environment="PATH=/var/www/nxsocial/venv/bin"
EnvironmentFile=/var/www/nxsocial/Socialmedia/.env

ExecStart=/var/www/nxsocial/venv/bin/celery \
    -A hub worker \
    --loglevel=info \
    --concurrency=2 \
    --logfile=/var/www/nxsocial/Socialmedia/logs/celery_worker.log \
    --pidfile=/tmp/celery_worker.pid \
    --detach

ExecStop=/bin/sh -c \
    'if [ -f /tmp/celery_worker.pid ]; then \
        kill -TERM $(cat /tmp/celery_worker.pid); \
    fi'

ExecReload=/bin/sh -c \
    'if [ -f /tmp/celery_worker.pid ]; then \
        kill -HUP $(cat /tmp/celery_worker.pid); \
    fi'

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celery-worker

[Install]
WantedBy=multi-user.target
