# celery-beat.service
[Unit]
Description=Celery Beat Scheduler for NxSocial
Documentation=https://docs.celeryq.dev
After=network.target redis.service celery-worker.service

[Service]
Type=forking
User=www-data
Group=www-data

# ── IMPORTANT: Update these paths to match your server setup ──
WorkingDirectory=/var/www/nxsocial/Socialmedia
Environment="PATH=/var/www/nxsocial/venv/bin"
EnvironmentFile=/var/www/nxsocial/Socialmedia/.env

ExecStart=/var/www/nxsocial/venv/bin/celery \
    -A hub beat \
    --loglevel=info \
    --logfile=/var/www/nxsocial/Socialmedia/logs/celery_beat.log \
    --pidfile=/tmp/celery_beat.pid \
    --schedule=/tmp/celerybeat-schedule \
    --detach

ExecStop=/bin/sh -c \
    'if [ -f /tmp/celery_beat.pid ]; then \
        kill -TERM $(cat /tmp/celery_beat.pid); \
    fi'

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celery-beat

[Install]
WantedBy=multi-user.target
